# Scripts/cron/crontab
#RPOps
0 7 * * * /home/sustadmin/PDI/jobs/startWithCarte.sh w_RPOps_AWS_Sync.kjb